# check_arp.service
[Unit]
Description  = Run check arp script
After=default.target

[Service]
ExecStart=/etc/check_arp.sh

[Install]
WantedBy=default.target
